# agent/zoraxy-agent.service
[Unit]
Description=Zoraxy Agent - Remote VPS management for ZoraxyHub
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/zoraxy-agent
ExecStart=/usr/local/bin/bun run zoraxy-agent.ts
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

Environment=NODE_ENV=production

StandardOutput=journal
StandardError=journal
SyslogIdentifier=zoraxy-agent

[Install]
WantedBy=multi-user.target
